# supervisord.conf
[supervisord]
nodaemon=true

[program:jupyterhub]
command=jupyterhub -f /srv/jupyterhub/jupyterhub_config.py

[program:nslcd]
command=nslcd -d
