# notebook-image/files/supervisord_notebook.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/log

[program:jupyter notebook]
command=sh /srv/singleuser/systemuser.sh

[program:nslcd]
command=nslcd -d
