# configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/log

[program:jupyterhub]
command=jupyterhub -f /srv/jupyterhub/configs/jupyterhub_config.py

[program:nslcd]
command=nslcd -d
